# com.openclaw.guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.guardian</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/haoc/.openclaw/workspace/openclaw_guardian.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/haoc/.openclaw/workspace/logs/guardian_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/haoc/.openclaw/workspace/logs/guardian_stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/haoc/.openclaw/workspace</string>
</dict>
</plist>
